# deploy/yumi-health.timer
[Unit]
Description=Yumi health check cada 5 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
